# deploy/drop-ins/chad-trends-refresh.service
[Unit]
Description=CHAD Google Trends Refresh
After=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/chad_finale
Environment=PYTHONPATH=/home/ubuntu/chad_finale
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/home/ubuntu/chad_finale/venv/bin/python3 -m chad.ops.trends_refresh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chad-trends-refresh

[Install]
WantedBy=multi-user.target
