# ops/systemd/chad-ibkr-collector.service.d/10-timeout-guards.conf
# NEW-GAP-043 / Box-021 — finite timeout guardrails for chad-ibkr-collector.
#
# Pre-fix audit (2026-05-19T21:30:00Z) observed:
#   - chad-ibkr-collector.service:
#       Type=oneshot
#       TimeoutStartUSec=infinity      <- unbounded
#       RuntimeMaxUSec=infinity        <- unbounded
#       ActiveState=activating, SubState=start
#       MainPID=3728456, ExecMainStartTimestamp=2026-05-19 12:30:32 UTC
#       (i.e. stuck in `activating/start` for ~9h18m)
#   - chad-ibkr-collector.timer fires every 120s but cannot fire again
#     until the in-flight run completes, so 9h of timer ticks were lost.
#
# These drop-in directives establish finite bounds at the systemd layer
# while a companion code-level SIGALRM guard
# (chad/portfolio/ibkr_portfolio_collector_v2.py::install_wall_clock_guard)
# provides defense-in-depth even on stripped-down/test deployments.
#
# Why these numbers:
#   - Healthy collector runs complete in 2-5 seconds (connect + accountSummary
#     + positions + atomic_write). 60s start budget is ~12-30x the expected
#     wall clock, comfortably absorbs gateway transient stalls but cuts off
#     anything beyond reasonable.
#   - RuntimeMaxSec=90 caps total runtime including the post-start window
#     (oneshot exit is considered "still starting" until the process returns).
#   - TimeoutStopSec=30: SIGTERM → 30s grace → SIGKILL. Default is 90s
#     which is excessive for a process that should already have wall-clock
#     self-terminated via SIGALRM.
#
# Channel-1 operator deployment:
#   sudo mkdir -p /etc/systemd/system/chad-ibkr-collector.service.d
#   sudo cp ops/systemd/chad-ibkr-collector.service.d/10-timeout-guards.conf \
#           /etc/systemd/system/chad-ibkr-collector.service.d/
#   sudo systemctl daemon-reload
#   sudo systemctl reset-failed chad-ibkr-collector.service
#   # If the in-flight stuck PID survives, send SIGTERM:
#   #   sudo systemctl stop chad-ibkr-collector.service
#   sudo systemctl status chad-ibkr-collector.service
#
# Verification (after Channel-1):
#   systemctl show chad-ibkr-collector.service \
#     -p Type,TimeoutStartUSec,TimeoutStopUSec,RuntimeMaxUSec
#   Expected:
#     Type=oneshot
#     TimeoutStartUSec=1min      (60s)
#     TimeoutStopUSec=30s
#     RuntimeMaxUSec=1min 30s    (90s)

[Service]
TimeoutStartSec=60
TimeoutStopSec=30
RuntimeMaxSec=90
