# ops/systemd/chad-portfolio-merge.service
[Unit]
Description=CHAD Portfolio Snapshot Merger (IBKR + Kraken)
After=network-online.target chad-ibkr-collector.service chad-kraken-collector.service
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/chad_finale
Environment=PYTHONUNBUFFERED=1
ExecStart=/bin/bash -lc 'cd "/home/ubuntu/chad_finale" && exec "/home/ubuntu/CHAD FINALE/venv/bin/python3" -m chad.portfolio.merge_kraken_into_snapshot'

[Install]
WantedBy=multi-user.target
