# ops/systemd/chad-ibkr-watchdog.service
[Unit]
Description=CHAD IBKR Watchdog (Phase 9B Telegram alerts)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/chad_finale

Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/ubuntu/chad_finale
EnvironmentFile=-/etc/chad/telegram.env

# Optional overrides (safe defaults)
Environment=IBKR_WATCHDOG_MAX_FAILURES=1
Environment=IBKR_WATCHDOG_MAX_AGE_SECONDS=360
Environment=TELEGRAM_NOTIFY_DEDUPE_TTL_SECONDS=900

ExecStart=/bin/bash -lc 'cd "/home/ubuntu/chad_finale" && source "/home/ubuntu/chad_finale/venv/bin/activate" && exec python -m chad.core.ibkr_watchdog'

Nice=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=false
ReadWritePaths=/home/ubuntu/chad_finale/runtime

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
