# ops/systemd/chad-clean-soak-evaluator.timer
[Unit]
Description=CHAD Clean-Soak Evaluator Timer (every 5 minutes)

[Timer]
OnBootSec=120
OnUnitActiveSec=300
AccuracySec=30s
Persistent=true
Unit=chad-clean-soak-evaluator.service

[Install]
WantedBy=timers.target
