# ops/systemd_templates/chad-live-loop.service.d/90-fd-limit.conf
# CHAD live_loop FD safety guard.
#
# 2026-05-30 audit found the live loop running with:
#   soft=1024 hard=524288
# and eventually exhausting the soft limit, causing Errno 24 failures across
# bars, runtime JSON files, VIX, Kraken prices, dynamic caps, and guard files.
#
# This raises the active process FD limit. This is a mitigation/safety guard;
# the FD growth/root cause should still be audited separately.
[Service]
LimitNOFILE=65536
