# ops/systemd_templates/chad-live-loop.service.d/91-disable-futures-exec.conf
# Emergency safety guard — 2026-05-30
# Stops repeated futures order submission loop for MCL/M6E while code-level
# cooldown/rejection handling is patched.
#
# This should prevent live_loop from submitting FUT intents to IBKR.
[Service]
Environment=CHAD_DISABLE_FUTURES_EXECUTION=1
Environment=CHAD_DISABLE_FUTURES=1
Environment=CHAD_FUTURES_EXECUTION_ENABLED=0
