# deploy/chad-kraken-futures-intel-refresh.timer
[Unit]
Description=CHAD Kraken Futures public intelligence refresh — every 5 minutes

[Timer]
OnBootSec=150
OnUnitActiveSec=300
AccuracySec=15s
Persistent=true

[Install]
WantedBy=timers.target
