# deploy/chad-setup-expectancy.service
[Unit]
Description=CHAD setup family expectancy updater — computes per-family P&L metrics for alpha_intraday_micro
After=network.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/chad_finale
Environment=PYTHONPATH=/home/ubuntu/chad_finale
EnvironmentFile=-/etc/chad/chad.env
ExecStart=/home/ubuntu/chad_finale/venv/bin/python3 /home/ubuntu/chad_finale/ops/update_setup_family_expectancy.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chad-setup-expectancy

[Install]
WantedBy=multi-user.target
