# deploy/chad-symbol-blocker.service
[Unit]
Description=CHAD Symbol Performance Blocker (oneshot)
After=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/chad_finale
Environment=PYTHONPATH=/home/ubuntu/chad_finale
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/home/ubuntu/chad_finale/venv/bin/python -m chad.risk.symbol_performance_blocker
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chad-symbol-blocker
