# deploy/chad-rs-refresh.service
[Unit]
Description=CHAD relative strength publisher
After=network.target

[Service]
Type=oneshot
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/chad_finale
Environment=PYTHONPATH=/home/ubuntu/chad_finale
ExecStart=/home/ubuntu/chad_finale/venv/bin/python3 -m chad.market_data.relative_strength_publisher
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chad-rs-refresh

[Install]
WantedBy=multi-user.target
