# deploy/chad-short-interest-refresh.service
[Unit]
Description=CHAD Short Interest Refresh
After=network-online.target

[Service]
Type=oneshot
User=ubuntu
WorkingDirectory=/home/ubuntu/chad_finale
Environment=PYTHONPATH=/home/ubuntu/chad_finale
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/home/ubuntu/chad_finale/venv/bin/python3 -m chad.ops.short_interest_refresh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=chad-short-interest-refresh

[Install]
WantedBy=multi-user.target
