# ops/systemd/chad-clean-soak-evaluator.service.d/10-onfailure.conf
# P0A-A2 — route this oneshot's failures to the CHAD service-failure alert
# handler so a silent failure reaches the operator (Telegram + a
# reports/service_failures/ artifact) instead of going unnoticed.
#
# Motivating gap (census R43): chad-ibkr-daily-bars-refresh ran
# Active=failed (exit-1) for ~17h with NO OnFailure=, and downstream read its
# stale data/bars/1d/ output as fresh. Every failure-prone oneshot must alert.
#
# %N is the failed unit name WITHOUT its `.service` suffix. The template
# chad-service-alert@.service re-appends `.service` in its ExecStart
# (`--failed-unit %i.service`), so `%N` yields the correct instance name
# chad-service-alert@<unit>.service. Do NOT use `%n` here — it retains the
# `.service` suffix and would double it (chad-...service.service).
#
# Channel-2 design artifact only. Install via ops/systemd/INSTALL_onfailure.md
# (cp + `systemctl daemon-reload`). Never edit installed unit files directly
# (CLAUDE.md governance rule #6).
[Unit]
OnFailure=chad-service-alert@%N.service
